# target/darknode/provisions/darknode-updater.service
[Unit]
Description=Republic Protocol's Darknode Automatic Updater
After=network.target

[Service]
ExecStart=/bin/bash $HOME/.darknode/updater.sh
Restart=on-failure
StartLimitBurst=0

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
# SIGHUP gives parity time to exit cleanly before SIGKILL (default 90s)
KillSignal=SIGHUP

[Install]
WantedBy=default.target