# target/darknode/provisions/logstash.service
[Unit]
Description=Logstash Daemon
After=network.target

[Service]
# run as root, set base_path in config.toml
ExecStart= $HOME/logstash-6.2.2/bin/logstash -f $HOME/logstash-6.2.2/darknode.conf
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
# SIGHUP gives parity time to exit cleanly before SIGKILL (default 90s)
KillSignal=SIGHUP

[Install]
WantedBy=default.target