# darknode/provisions/rc.local
go get -u github.com/republicprotocol/republic-go/cmd/dark-node
rm darknode.log
sudo systemctl restart darknode.service
exit 0